# programs.conf
[program:rabbitmq]
command=rabbitmq-server
autostart=true
stopasgroup=true
killasgroup=true

[program:confd]
command=confd --backend rancher --prefix /2015-07-25 %(ENV_CONFD_ARGS)s
autostart=true
